# .config/systemd/user/redshift.service.d/override.conf
[Unit]
After=network-online.target

[Service]
ExecStart=
ExecStart=redshift -l 33:-96 -v
StartLimitIntervalSec=0